# node-setup/traffic-exporter.service
[Unit]
Description=Traffic exporter for Grafana Network Observability Platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/grafana-netobs/traffic-exporter.env
WorkingDirectory=/opt/grafana-network-observability-platform/exporters/traffic_exporter
ExecStart=/opt/grafana-netobs-venv/bin/python /opt/grafana-network-observability-platform/exporters/traffic_exporter/main.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target